# hello-world.service.in
[Unit]
Description=Hello World

[Service]
Type=simple
ExecStart=@EXE_PATH@

[Install]
WantedBy=multi-user.target